# demo-service/supervisord.conf
[supervisord]
nodaemon=true

[program:demo-service]
command=/app/demo-service